# src/lib/systemd/system/top3.timer
[Unit]
Description=Update top3 every minutes

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy = timers.target
