# rpimonitor/web/addons/top3/top3.cron
# 
* * * * * root cd /usr/share/rpimonitor/web/addons/top3; ./top3 > top3.html

