# init/sysv/rpimonitor
#!/bin/bash
#
# This starts and stops RPI-Monitor
#
# chkconfig: 2345 12 88
# description: rpimonitor is a slef monitoring tool for RaspberryPi
# processname: rpimonitord
# pidfile: /var/run/rpimonitord.pid
### BEGIN INIT INFO
# Provides: $XXXXXXXXXXXX
### END INIT INFO

# Source function library.
. /lib/lsb/init-functions

DAEMON="/usr/local/rpimonitor/rpimonitord"
PIDFILE="/var/run/rpimonitord.pid"

[ -x $PROGRAM ] || exit 0



start() {
    log_daemon_msg "Starting RPi-Monitor" "rpimonitord"
    $DAEMON -b
    status=$?
    log_end_msg $status
}

stop() {
    log_daemon_msg "Stopping RPi-Monitor" "rpimonitord"
    for pid in $(cat $PIDFILE); do
      kill -9 $pid
    done
    status=$?
    log_end_msg $status
}

restart() {
    echo -n "Restarting RPi-Minitor"
    stop
    sleep 2
    start
}

case "$1" in
    start)
    	start
    ;;
    stop)
    	stop
    ;;
    restart)
    	restart
    ;;
    *)
    	echo "Usage: $0 {start|stop|restart}"
    ;;
esac

exit 0
